# demo/docker/vulnerable/supervisord.conf
[inet_http_server]
port=0.0.0.0:9001

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
priority=10
autostart=true
autorestart=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=20
autostart=true
autorestart=true
user=root

[program:queue-worker]
command=/bin/sh -c "while true; do php /var/www/html/artisan queue:work --tries=0 --timeout=0 || true; sleep 5; done"
priority=30
autostart=true
autorestart=true
user=root

[program:scheduler]
command=/bin/sh -c "while true; do php /var/www/html/artisan schedule:run || true; sleep 60; done"
priority=40
autostart=true
autorestart=true
user=root
